# package/root.pacman/etc/systemd/system/rtk-hciattach.service
[Unit]
Description=Realtek H5 bluetooth support

[Service]
ExecStart=/usr/local/sbin/pinebook_start_bt.sh

[Install]
WantedBy=multi-user.target
